# data/wikioracle.service
# WikiOracle — systemd service
# Install: sudo cp wikioracle.service /etc/systemd/system/
#          sudo systemctl daemon-reload && sudo systemctl enable wikioracle
#
# The shim runs in stateless mode behind the WordPress/Apache reverse proxy.
# Apache ProxyPass /client/ https://127.0.0.1:8787/client/
# Apache ProxyPassReverse /client/ https://127.0.0.1:8787/client/
#
# IMPORTANT: Only /client should be proxied. /nanochat is NOT proxied — the
# NanoChat inference endpoint is internal-only (127.0.0.1:8000). The shim
# calls it directly via WIKIORACLE_BASE_URL.

[Unit]
Description=WikiOracle Chat Shim (stateless Flask)
After=network.target

[Service]
Type=simple
User=bitnami
WorkingDirectory=/opt/bitnami/wordpress/files/WikiOracle.org/client
ExecStart=/opt/bitnami/wordpress/files/WikiOracle.org/client/.venv/bin/python3 \
    bin/wikioracle.py --stateless --url-prefix /client
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=WIKIORACLE_BIND_HOST=127.0.0.1
Environment=WIKIORACLE_BIND_PORT=8787
Environment=WIKIORACLE_BASE_URL=http://127.0.0.1:8000
Environment=WIKIORACLE_API_PATH=/chat/completions
Environment=WIKIORACLE_STATE_FILE=/opt/bitnami/wordpress/files/WikiOracle.org/client/state.xml
# Bearer-token auth on all endpoints (except /health).
# Generate with: python3 -c "import secrets; print(secrets.token_urlsafe(32))"
Environment=WIKIORACLE_API_TOKEN=CHANGEME

[Install]
WantedBy=multi-user.target
